# systemd/salesup-research-worker.service
[Unit]
Description=SalesUp Deep Research Worker
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/salesup-team-bot
EnvironmentFile=/opt/salesup-team-bot/.env
Environment="PATH=/opt/salesup-team-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/salesup-team-bot/venv/bin/python3 -m research_worker
Restart=on-failure
RestartSec=10
MemoryMax=1536M
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
